# echo/host/etc/systemd/system/certbot.service.d/override.conf
[Service]
ExecStartPost=systemctl reload nginx.service
