# jitsi/jitok/container/etc/systemd/system/jitok.service
[Unit]
Description=Jitsi Token API Service
After=network-online.target

[Service]
User=jitok
Group=jitok
ExecStart=/home/jitok/api/jitok.sh
KillMode=control-group
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
